# startkde/plasma-base.target
[Unit]
Description=KDE Plasma Workspace
Requires=plasma-core.target
