# startkde/plasma-core.target
[Unit]
Description=KDE Plasma Workspace Core
Wants=plasmashell.service kwin_x11.service ksmserver.service kcminit.service kded.service
