# startkde/plasma-core.target
[Unit]
Description=KDE Plasma Workspace
Requires=plasmashell.service kwin.service ksmserver.service kcminit.service kded.service
