# startkde/plasma-workspace.target
[Unit]
Description=KDE Plasma Workspace
Requires=plasma-core.target
Wants=ksplash-ready.service
