# ksmserver/plasma-ksmserver.service.in
[Unit]
Description=KDE Session Management Server
Wants=plasma-kcminit.service
PartOf=graphical-session.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/ksmserver
BusName=org.kde.ksmserver
Slice=session.slice

[Install]
WantedBy=plasma-core.target
