# ksmserver/plasma-restoresession.service.in
[Unit]
Description=KDE Session Management Server
After=graphical-session.target
RefuseManualStart=yes

[Service]
Type=oneshot
ExecStart=-@QtBinariesDir@/qdbus org.kde.ksmserver /KSMServer org.kde.KSMServerInterface.restoreSession
Slice=session.slice
