# shell/plasma-plasmashell.service.in
[Unit]
Description=KDE Plasma Workspace
After=plasma-ksmserver.service plasma-kcminit.service
PartOf=graphical-session.target
StartLimitIntervalSec=60s
StartLimitBurst=3

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/plasmashell --no-respawn
Restart=on-failure
Type=dbus
BusName=org.kde.plasmashell
Slice=session.slice
TimeoutSec=40sec

[Install]
WantedBy=plasma-core.target
