# xembed-sni-proxy/plasma-xembedsniproxy.service.in
[Unit]
Description=Handle legacy xembed system tray icons
PartOf=graphical-session.target
After=plasma-core.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/xembedsniproxy
Restart=on-failure
Type=simple
Slice=background.slice
TimeoutSec=5sec
