# krunner/plasma-krunner.service.in
[Unit]
Description=KRunner
PartOf=graphical-session.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/krunner
Type=dbus
BusName=org.kde.krunner
TimeoutSec=5sec
Slice=background.slice
# Disable restart as we're dbus activated anyway
Restart=no
