# startkde/kcminit/plasma-kcminit-phase1.service.in
[Unit]
Description=KDE Configuration Module Initialization (Phase 1)
Requires=kcminit.service
After=plasma-kcminit.service plasma-kded.service

[Service]
Type=oneshot
ExecStart=@QtBinariesDir@/qdbus org.kde.kcminit /kcminit org.kde.KCMInit.runPhase1
Slice=session.slice
