# runners/baloo/plasma-baloorunner.service.in
[Unit]
Description=KRunner provider for baloo file indexer
PartOf=graphical-session.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/baloorunner
ExecCondition=@KDE_INSTALL_FULL_BINDIR@/kde-systemd-start-condition --condition "baloofilerc:Basic Settings:Indexing-Enabled:true"
Type=dbus
BusName=org.kde.runners.baloo
TimeoutSec=5sec
Slice=background.slice
Restart=no #as we're dbus activated anyway
