# startkde/systemd/plasma-login-wayland.target
[Unit]
Requires=plasma-kwin_wayland.service
BindsTo=plasma-kwin_wayland.service
Requires=plasma-login.service
BindsTo=plasma-login.service
