# startkde/systemd/plasma-ksplash-ready.service.in
[Unit]
Description=KSplash "ready" Stage
Wants=plasma-core.target
After=plasma-core.target
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=dbus-send --session  --reply-timeout=1 --dest=org.kde.KSplash /KSplash org.kde.KSplash.setStage string:ready
Slice=session.slice
