# startkde/systemd/plasma-sentinel.service.in
# The idea to bind this to plasma-workspace.target
# That part is disabled as it causes issues with the hot-upgrade path

[Unit]
Description=Close Plasma session when session leader quits
PartOf=graphical-session.target
RefuseManualStart=yes
RefuseManualStop=yes

[Service]
Type=exec
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/plasma_waitforname org.kde.startplasma --unregistered --timeout 0
Slice=session.slice
